# rds-proxy/proxy/bin/RdsProxy.service
# cp this file to /etc/systemd/system/
# and change 'ExecStart' and 'ExecStop' to correct

[Unit]
Description=RDS Proxy Service

[Service]
Type=simple
ExecStart=/opt/rds/rdsproxy/bin/StartProxy.sh
ExecStop=/opt/rds/rdsproxy/bin/StopProxy.sh
PrivateTmp=false

TimeoutSec=10

[Install]
WantedBy=multi-user.target